# install/etc/profile.d/pvzdweb.sh
source /opt/venv/pvzdweb/bin/activate
export PROJ_HOME=/opt/PVZDweb
source /opt/PVZDweb/bin/setenv.sh