# install/etc/profile.d/pvzdweb.sh
source /opt/venv/pvzdweb/bin/activate
source /opt/PVZDweb/bin/setenv.sh
